# api/crontab
#Mins   Hours   Days    Months    Weekday   Command
*       */2       *       *         *         npm run --prefix /app import-data >> /app/cronjob.log 2>&1
